# resources/docs/platform/nginx-java/docker/config/supervisor/conf.d/java.conf.sample
[program:java]
command=java -Xms256m -Xmx512m -jar /var/www/app.jar
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0